# easylinker-proxy-service/src/main/resources/etc/server/supervisor_config.conf
#supervisor 是一个进程管理工具，可以很好的监视管理进程
#关于此软软件请上网查些资料
#下面是EasyLinker 再supervisor下的配置
#yum install supervisor
#systemctl enable supervisord.service
#启动 systemctl start supervisord.service
#关闭 systemctl stop supervisord.service
#重启 systemctl restart supervisord.service
#
#
[program:easylinker]
command=nohup java -jar EasyLinkerApplication-0.0.1-SNAPSHOT.jar --spring.config.location = application.properties &
directory= /root/easylinker/#目录自定
user=root#root权限
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/easylinker.log
loglevel=info